# Assigments/tvs-2324-1-leic51n-coursework-g03n1/cw3/ex2/tvsapp.service
[Unit]
Description=daemon for tvsapp service
Requires=nginx.service
Requires=tvsapp.socket

[Service]
ExecStart=/opt/isel/tvs/tvsctld/bin/tvsapp_server

[Install]
WantedBy=multi-user.target