# Assigments/tvs-2324-1-leic51n-coursework-g03n1/cw3/ex2/tvsapp.socket
[Unit]
Description=TVS daemon activation socket

[Socket]
ListenStream=/run/isel/tvsctld/request

[Install]
WantedBy=sockets.target