# Assigments/tvs-2223-1-leic52d-g01/se3/ex5/tvsctl.socket
[Unit]
Description=tvsctl socket

[Socket]
ListenStream=/run/isel/tvsctl/request

[Install]
WantedBy=sockets.target
