# Assigments/tvs-2223-1-leic52d-g01/se3/ex5/tvsctl.service
[Unit]
Description=TVS Control Application
Requires=nginx.service
Requires=tvsctl.socket

[Service]
ExecStart=/opt/isel/tvs/tvsctl/bin/tvsctl
Restart=no

[Install]
Also=tvsapp@.service
WantedBy=multi-user.target

